# router/files/etc/uci-defaults/dhcp-relay
# Trash ISP (Vodafone Kabel) only hands out /64s.
# So, relay mode.

set -uexo pipefail

uci -q batch <<EOI
set dhcp.lan.dhcpv6=relay
set dhcp.lan.ra=relay
set dhcp.lan.ndp=relay

set dhcp.wan6=dhcp
set dhcp.wan6.dhcpv6=relay
set dhcp.wan6.ra=relay
set dhcp.wan6.ndp=relay
set dhcp.wan6.master=1
set dhcp.wan6.interface=wan6

commit dhcp
EOI

service dnsmasq restart || true
service odhcpd restart || true
