# router/files/etc/uci-defaults/fp2-usb
set -uexo pipefail

uci -q batch << EOI
set network.fp2dev=device
set network.fp2dev.name='usb0'

set network.fp2=interface
set network.fp2.proto='dhcp'
set network.fp2.device='usb0'
set network.fp2.defaultroute='0'
set network.fp2.peerdns='0'
set network.fp2.delegate='0'
commit network
EOI
