# router/files/etc/uci-defaults/pxe
# Configures dnsmasq for PXE boot.
# This config assumes that:
# * external storage, mounted at `/mnt/storage`, has a tftp/ directory
#   containing iPXE undionly.kpxe + ipxe.efi.
# * /mnt/storage/coreos is bind mounted to `/www/coreos`
# * coreos directory contains a boot.ipxe script.

# iPXE is chainloaded, and then it is instructed to run boot.ipxe script,
# which is loaded from the uHTTPd server (the one that powers LUCI).

# TODO: update reference to source tree that prepares this storage location,
# when it exists.

set -uexo pipefail

uci -q batch <<EOI
set dhcp.@dnsmasq[0].enable_tftp='1'
set dhcp.@dnsmasq[0].tftp_root='/mnt/storage/coreos'
commit dhcp
EOI

cat > /etc/dnsmasq.conf <<HERE
dhcp-boot=tag:home-cluster,ipxe.efi
HERE

service dnsmasq restart || true
