# router/files/etc/uci-defaults/restricted-network
# Configures an SSID which has very limited connectivity.
# Devices in this network can talk to:
#  * DNS/DHCP
#  * samcday-headscale.fly.dev
#  * UDP port 41641 on remote hosts 
#  * That's it.
# The idea here is to provide enough external access to connect an untrusted device to the tailnet.
# Everything else (including Internet egress) should be handled from inside the tailnet.

set -uexo pipefail

# TODO: move this elsewhere
exit 0

# adapted from https://openwrt.org/docs/guide-user/network/wifi/guestwifi/guest-wlan

# Network interface + bridge
uci -q batch << EOI
set network.restricted_dev="device"
set network.restricted_dev.type="bridge"
set network.restricted_dev.name="br-restricted"

set network.restricted="interface"
set network.restricted.proto="static"
set network.restricted.device="br-restricted"
set network.restricted.ipaddr="${IPADDR_RESTRICTED}"
set network.restricted.netmask="255.255.255.0"
commit network
EOI

# DHCP server + ipset
uci -q batch << EOI
set dhcp.restricted="dhcp"
set dhcp.restricted.interface="restricted"
set dhcp.restricted.start="100"
set dhcp.restricted.limit="150"
set dhcp.restricted.leasetime="1h"
set dhcp.restricted.netmask="255.255.255.0"

set dhcp.restricted-ipset="ipset"
add_list dhcp.restricted-ipset.name="restricted-ipset"
add_list dhcp.restricted-ipset.name="restricted-ipset6"
add_list dhcp.restricted-ipset.domain="samcday-headscale.fly.dev"
commit dhcp
EOI
/etc/init.d/dnsmasq restart || true

# Wireless SSID
uci -q batch << EOI
set wireless.restricted="wifi-iface"
set wireless.restricted.device="$(uci get wireless.@wifi-iface[1].device)" # iface 1 is 5ghz.
set wireless.restricted.mode="ap"
set wireless.restricted.network="restricted"
set wireless.restricted.ssid="${SSID}-restricted"
set wireless.restricted.encryption="none"
commit wireless
EOI
wifi reload

# Firewall config
uci -q batch << EOI
set firewall.restricted="zone"
set firewall.restricted.name="restricted"
set firewall.restricted.network="restricted"
set firewall.restricted.input="ACCEPT"
set firewall.restricted.output="ACCEPT"
set firewall.restricted.forward="REJECT"

set firewall.restricted_dns="rule"
set firewall.restricted_dns.name="restricted_dns"
set firewall.restricted_dns.src="restricted"
set firewall.restricted_dns.dest_port="53"
set firewall.restricted_dns.proto="tcp udp"
set firewall.restricted_dns.target="ACCEPT"

set firewall.restricted_dhcp="rule"
set firewall.restricted_dhcp.name="restricted_dhcp"
set firewall.restricted_dhcp.src="restricted"
set firewall.restricted_dhcp.src_port="68"
set firewall.restricted_dhcp.dest_port="67"
set firewall.restricted_dhcp.proto="udp"
set firewall.restricted_dhcp.family="ipv4"
set firewall.restricted_dhcp.target="ACCEPT"

set firewall.restricted_tailnet_cp="rule"
set firewall.restricted_tailnet_cp.name="restricted_tailnet_cp"
set firewall.restricted_tailnet_cp.src="restricted"
set firewall.restricted_tailnet_cp.dest_port="443"
set firewall.restricted_tailnet_cp.dest="wan"
set firewall.restricted_tailnet_cp.dest_ip="137.66.18.54"
set firewall.restricted_tailnet_cp.proto="tcp"
set firewall.restricted_tailnet_cp.target="ACCEPT"

set firewall.restricted_tailscale="rule"
set firewall.restricted_tailscale.name="restricted_tailscale"
set firewall.restricted_tailscale.src="restricted"
set firewall.restricted_tailscale.dest_port="41641"
set firewall.restricted_tailscale.dest="*"
set firewall.restricted_tailscale.proto="udp"
set firewall.restricted_tailscale.target="ACCEPT"

commit firewall
EOI
/etc/init.d/firewall restart
