# router/files/etc/uci-defaults/root-pw
set -uexo pipefail

passwd root <<EOP
${ROOT_PW}
${ROOT_PW}
EOP
