# router/files/etc/uci-defaults/system
set -uexo pipefail

uci -q batch << EOI
set system.@system[0].zonename='Europe/Berlin'
set system.@system[0].hostname='${HOSTNAME}'
commit system
set network.lan.ipaddr='${IPADDR}'
commit network
EOI

ntpd -q -p 0.openwrt.pool.ntp.org
