# devices/router/files/etc/uci-defaults/00-mesh11sd
set -uexo pipefail

uci -q batch <<EOI
set mesh11sd.setup.auto_config='1'
set mesh11sd.setup.auto_mesh_id='sam-home-cluster'
set mesh11sd.setup.auto_mesh_band='5g'
set mesh11sd.setup.mesh_gate_encryption='1'
set mesh11sd.setup.ssid_suffix_enable='0'
EOI
