# devices/router/files/etc/uci-defaults/00-wifi
set -uexo pipefail

uci -q batch << EOI
set wireless.@wifi-device[0].disabled='0'
set wireless.@wifi-device[1].disabled='0'

set wireless.@wifi-iface[0].ssid='samnet'
set wireless.@wifi-iface[0].encryption=sae
set wireless.@wifi-iface[0].ieee80211r='1'
set wireless.@wifi-iface[0].mobility_domain='babe'

set wireless.@wifi-iface[1].ssid='samnet'
set wireless.@wifi-iface[1].encryption=sae
set wireless.@wifi-iface[1].ieee80211r='1'
set wireless.@wifi-iface[1].mobility_domain='babe'
EOI

# The commit + wifi reload happens in 99-wifi-password
