# devices/router/files/etc/uci-defaults/00-wifi
set -uexo pipefail

uci -q batch << EOI
set wireless.@wifi-device[0].country='DE'
set wireless.@wifi-device[1].country='DE'
commit wireless
EOI
