# devices/router/files/etc/uci-defaults/dns-rebind
# Tailnet-internal services are currently reached via the service IPs
# (172.31.0.0/16) directly. OpenWRT rebind protection freaks out about this.

set -uexo pipefail

uci -q batch <<EOI
add_list dhcp.cfg01411c.rebind_domain='samcday.com'
commit dhcp
EOI

service dnsmasq restart || true
service odhcpd restart || true
