# devices/router/files/etc/uci-defaults/mwan3
set -uexo pipefail

# static mwan3 defaults to:
# use (only) wired wan port if available
# use (only) wwan if available
# use (loadbalanced..?) usb tethers/dongles if available

uci delete mwan3.wan_m1_w3 || true
uci delete mwan3.wan_m2_w3 || true
uci delete mwan3.wanb_m1_w2 || true
uci delete mwan3.wanb_m1_w3 || true
uci delete mwan3.wanb_m2_w2 || true
uci delete mwan3.wan6_m1_w3 || true
uci delete mwan3.wan6_m2_w3 || true
uci delete mwan3.wanb6_m1_w2 || true
uci delete mwan3.wanb6_m1_w3 || true
uci delete mwan3.wanb6_m2_w2 || true

uci delete mwan3.wan6 || true
uci delete mwan3.wan || true
uci delete mwan3.wanb || true
uci delete mwan3.wanb6 || true

for policy in wan_only wanb_only balanced wan_wanb wanb_wan; do
    uci delete mwan3.${policy} || true
done

for rule in https default_rule_v4 default_rule_v6; do
    uci delete mwan3.${rule} || true
done

uci delete mwan3.all || true
uci set mwan3.all=policy
uci set mwan3.all.last_resort='unreachable'

uci delete mwan3.default || true
uci set mwan3.default=rule
uci set mwan3.default.proto='all'
uci set mwan3.default.sticky='0'
uci set mwan3.default.use_policy='all'

for iface in wan wwan ipheth cdceth; do
    uci set mwan3.m_${iface}=member
    uci set mwan3.m_${iface}.interface="${iface}"
    uci set mwan3.m_${iface}.weight='1'
    uci add_list mwan3.all.use_member="m_${iface}"

    uci set mwan3.${iface}=interface
    uci set mwan3.${iface}.enabled='1'
    uci add_list mwan3.${iface}.track_ip=1.0.0.1
    uci add_list mwan3.${iface}.track_ip=1.1.1.1
    uci add_list mwan3.${iface}.track_ip=208.67.222.222
    uci add_list mwan3.${iface}.track_ip=208.67.220.220
    uci add_list mwan3.${iface}.track_ip=8.8.4.4
    uci add_list mwan3.${iface}.track_ip=8.8.8.8
    uci set mwan3.${iface}.family='ipv4'
    uci set mwan3.${iface}.reliability='2'
    uci set mwan3.${iface}.up='2'
done

uci set mwan3.m_wan.metric='1'
uci set mwan3.m_wwan.metric='5'
uci set mwan3.m_ipheth.metric='10'
uci set mwan3.m_cdceth.metric='10'

uci commit mwan3

service mwan3 reload || true
