# devices/router/files/etc/uci-defaults/prometheus-node-exporter-lua
set -uexo pipefail

uci -q batch << EOI
set prometheus-node-exporter-lua.main.listen_interface=*
commit prometheus-node-exporter-lua
EOI

service prometheus-node-exporter-lua restart || true
