# devices/router/files/etc/uci-defaults/system
set -uexo pipefail

uci -q batch << EOI
set system.@system[0].zonename='Europe/Berlin'
set system.@system[0].hostname='${HOSTNAME}'

set system.@system[0].log_ip='172.31.0.9'
set system.@system[0].log_port='514'
set system.@system[0].log_proto='udp'

commit system
set network.lan.ipaddr='${IPADDR}'
commit network
EOI

service log restart || true
ntpd -q -p 0.openwrt.pool.ntp.org || true

chsh -s /bin/bash || true
