# devices/router/files.enc/etc/uci-defaults/99-mesh11sd
{
	"data": "ENC[AES256_GCM,data:eJDaIbkRlnbT+cQGUT6VyGptJ7wB2MjEBRyVTAhh6ii1PllSNDzr/SwzdikTjETj5fil1pkuKzByS0213N+ePfSbpUI5liqcqXZC8eUP2B4WwJQ6nWbYha3TuwzLJYxVpjbARxnFqivuF9fUmGsR6xFqXIaOlkQLadHWPNE=,iv:1GAgLmK7Y3Hws6gscknFkWaw+75CbEq7CUNCSY30AjA=,tag:wjLQw+k4i/uDsQ9aTf5H5A==,type:str]",
	"sops": {
		"kms": null,
		"gcp_kms": null,
		"azure_kv": null,
		"hc_vault": null,
		"age": [
			{
				"recipient": "age10c9vvuvfkflc7zypu6zm8dtw0gdn028nlr3gslt35df8vdqrap5q36xav4",
				"enc": "-----BEGIN AGE ENCRYPTED FILE-----\nYWdlLWVuY3J5cHRpb24ub3JnL3YxCi0+IFgyNTUxOSAyNWZSNlZTY2tZNENCMDdm\nU0d0U3pqeGpNWHJ3Q1NaQ0lUUEhUNU1HYWs4CmZDbitRMFdxU3pEVXgyaXhheWF3\nUjZlclIxY2V5NEJoYTg3eEdXNGxMbE0KLS0tIFpoME1lQ2ZJSW1QbUJUQnNEVjB4\nbUFNMGZNeUNkN1VhYlhrK1hNTXVBVDAKeEgZgZhpJzoQU0xDAjJc0nXSKMgwEs/w\nuJkYFCYGRlQILw8K+4h8JothYuvNuKGqcx+z9tY0Df8oG8KqKm8pEA==\n-----END AGE ENCRYPTED FILE-----\n"
			}
		],
		"lastmodified": "2024-08-01T15:00:33Z",
		"mac": "ENC[AES256_GCM,data:ZR0u+PmQbRLU807/A/7kXaFpitypR0FyPxncGgjCS1nuvs7dRfpG0t3to6y/INjLv3TD90DYAxlApayVwo5LN58jzmBRo0MvzeTzWjpeBnCHuNjrbBkoQ0qFqdx8SeAD3V/TwVC9iK97WPzmMzrRubTRwBT4Z9bx4muxDWNhheo=,iv:ghoaqX+icC0HjNJuuExdYE5zog9mxhLFc8pvFfk7hvk=,tag:iZlWZG8G8DpNZZDqffMMhw==,type:str]",
		"pgp": null,
		"unencrypted_suffix": "_unencrypted",
		"version": "3.9.0"
	}
}