# boot/ignition-files/tree/etc/systemd/system/haproxy.service.d/kube-vip.conf
[Service]
ExecStartPre=-ip addr add 172.29.0.1/32 dev enp0s31f6
