# router/files/etc/rc.local
set -uexo pipefail

function bind_mount() {
    mkdir -p $1 $2
    mount --bind $1 $2
}

block mount

bind_mount /mnt/storage/coreos /www/coreos
bind_mount /mnt/storage/tang-db /usr/share/tang/db
bind_mount /mnt/storage/tailscale /etc/tailscale
service tailscale restart

if ! grep tangd < /etc/services; then
  echo "tangd 8888/tcp" >> /etc/services
  service xinetd restart
fi
