# node/usr/lib/systemd/system/rpm-ostreed.service.d/http-proxy.conf
[Service]
Environment="http_proxy=http://squid.samcday.com"
Environment="https_proxy=http://squid.samcday.com"
