# node/usr/lib/systemd/system/squid.service.d/ssldb-init.conf
[Service]
ExecStartPre=bash -c '[[ -f /var/spool/squid/ssl_db/index.txt ]] || /usr/lib64/squid/security_file_certgen -c -s /var/spool/squid/ssl_db -M 4MB'
