# router/files.enc/etc/uci-defaults/99-wifi-password
{
	"data": "ENC[AES256_GCM,data:qjSpqdq5G9LRxN7qjuhhqDFrGckO012nlP45nLD5btNGvog+cn69oap7Bd0PwMFKAHFeL0L3M4kiEgtOJyeGzYDoBjjpYyG5RvpSVGDeyeMvwMH+yUULCWGsbZ9dQiSpCb5T2sSN+Pys8KfDgJdA8+mzaDybia48ynkxX+2cJhboHOjmyRQpfy4KhCGGHZ3LH/0gpmx6ZnWJNRpD0jHKxoEtpMinw2XzDTaV0yHIYiVWBY/xFkjOImpZJuaFchyTXo5N6YOuOjyj0Edl6TCU5rwhAVqk1ncsLSRUcc/ieR9IEtstHEQxwUs=,iv:jNM8vNbV0n+jfYiXwKhoDQLIHNGtjSkQ8PTPKQJyoSw=,tag:CDrq/PMtm3ltRy72BwzwcQ==,type:str]",
	"sops": {
		"kms": null,
		"gcp_kms": null,
		"azure_kv": null,
		"hc_vault": null,
		"age": [
			{
				"recipient": "age10c9vvuvfkflc7zypu6zm8dtw0gdn028nlr3gslt35df8vdqrap5q36xav4",
				"enc": "-----BEGIN AGE ENCRYPTED FILE-----\nYWdlLWVuY3J5cHRpb24ub3JnL3YxCi0+IFgyNTUxOSBLeWxQVEZET0pkSHZ0TG1s\neVBYOUNXZTdrM2drT2d0bzA3L1JHZFNYR1UwCjcwdzI0L0lZZkhzcTRXd010bWt4\nMnNqNEZtMzFTRVBSZGdWbGhqaFdYRnMKLS0tIENZRDllTFM1OGU1UWNxR1NaRmpS\nNm9CVUtpOVpDL3lwZ0RLS2dqbUlCY3MK/pRQ9f3jVJxshI4fWV6k5tapWljVxhGN\n5O8hFLVzfxSkJvAvlSiGmzd8Zc8g3PGm7hel8ka+60SaYWtoPXGZqQ==\n-----END AGE ENCRYPTED FILE-----\n"
			}
		],
		"lastmodified": "2024-05-10T17:24:58Z",
		"mac": "ENC[AES256_GCM,data:8zHVYInbDRei8kYU5Ve3M7nqD0juTXACqoaA6CuHyilo4bfELvXYawUzxIhvzPvs3Lsmq943Krpx/5L43ifnjOD/uYkt1qWK63zKnaYR5+hipW71+KTmi51DqtGbXNtuwHn76GuIUs5SwyizcBHW1aJNlySdEWhWeuiZG3z7YsU=,iv:LIWQcZsvd04VR2jouOW3a51PEsE1J9fSfnI0jTnOutw=,tag:iJimD+SSZEcrRS4gpwMz3w==,type:str]",
		"pgp": null,
		"unencrypted_suffix": "_unencrypted",
		"version": "3.8.1"
	}
}