# devices/librem13/mkosi.extra/etc/systemd/system/tailscaled.service.d/logging.conf
[Service]
LogFilterPatterns=~Accept:
LogFilterPatterns=~portmapper: failed to get PCP mapping
LogFilterPatterns=~magicsock:.*now using
