# node/usr/lib/systemd/system/etc-kubernetes.mount
[Mount]
What=/var/lib/kubernetes
Where=/etc/kubernetes
Type=none
Options=bind

[Install]
WantedBy=multi-user.target
