# router/files/etc/hotplug.d/ntp/tailscale-up
set -uexo pipefail

[ "$ACTION" = "stratum" ] || exit 0

tailscale up \
  --login-server=https://samcday-headscale.fly.dev \
  --accept-routes \
  --auth-key=${TAILNET_AUTH_KEY} \
  --netfilter-mode=off \
 2>&1 | logger -st tailscale.up

rm /etc/hotplug.d/ntp/tailscale-up
